# .config/autostart/udiskie.desktop
[Desktop Entry]
Version=1.0
Name=Udiskie
Exec=udiskie -t
StartupNotify=false
Type=Application
