# .config/autostart/volumeicon.desktop
[Desktop Entry]
Version=1.0
Name=VolumeIcon
Exec=volumeicon
StartupNotify=false
Type=Application
